# client/douyu.conf
[program:app]
command=python /home/djf/douyu/douyu_client.py "/home/djf/douyu/movies1/"
autorstart=true 
stdout_logfile=/tmp/douyu.log