# DouyuApp.conf
[program:app]
command=python /opt/DouyuApp/main.py
autorstart=true 
stdout_logfile=/tmp/DouyuApp.log