# DouyuApp.conf
[program:app]
command=python /opt/app/DouyuApp/main.py
autorestart=true 
stdout_logfile=/tmp/DouyuApp.log